# etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
## NB: This service hangs in WSL, and is pretty meaningless anyways
[Service]
ExecStart=
ExecStart=/bin/true
